# deploy/schedulist-backend.service
[Unit]
Description=Schedulist Backend API
Documentation=https://github.com/YOUR_USERNAME/schedulist
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=schedulist
Group=schedulist
WorkingDirectory=/var/www/schedulist/schedulist
Environment="NODE_ENV=production"
Environment="PORT=5050"
ExecStart=/usr/bin/node /var/www/schedulist/schedulist/src/index.js
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/schedulist/backend.log
StandardError=append:/var/log/schedulist/backend-error.log

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/schedulist /var/log/schedulist /var/backups/schedulist

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target